# src/chassis/claude-code-channel/claude-code-channel@.service
# Claude Code Channel Instance Service Template
#
# Usage:
#   systemctl start claude-code-channel@Crossing-2d23
#   systemctl enable claude-code-channel@Crossing-2d23  (auto-start on boot)
#
# %i = instance ID (e.g., Crossing-2d23)
# %I = unescaped instance ID
#
# Prerequisites:
#   - claude-code-channel-setup.sh must have been run for this instance
#     (creates Unix user, .mcp.json, .claude.json bypass flags, etc.)
#   - .hacs-identity must exist in instance dir with channelPort allocated
#
# Install:
#   cp claude-code-channel@.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable claude-code-channel@<InstanceId>
#   systemctl start  claude-code-channel@<InstanceId>
#
# Design notes:
#   - Type=oneshot + RemainAfterExit=yes: ExecStart launches a detached tmux
#     session and exits cleanly. systemd treats the service as active until
#     ExecStop runs. This means **systemd is the launcher, not the
#     supervisor**: if claude-code dies inside tmux, systemd does NOT see
#     it. Restart=on-failure only fires if ExecStart itself fails.
#     In-flight liveness is monitored via the channel /health endpoint
#     (separate concern from systemd). A future /health watcher (cron or
#     Phase C /loop pulse) will provide active health monitoring; for
#     now, manual restart on /health failure is the recovery path.
#   - User=root: launch script internally uses `sudo -u %i` to drop
#     privileges before starting tmux. The chassis preserves Unix-username
#     case to match openfang and existing instance conventions, so %i maps
#     directly to the Unix user. Running as root lets the launch script
#     handle setup-vs-launch privilege orchestration uniformly.
#   - Logging via journal: `journalctl -u claude-code-channel@<id>`

[Unit]
Description=Claude Code Channel for HACS instance %i

# network-online.target (not network.target) — claude-code makes outbound
# calls at startup (Anthropic API, HACS registration, npm marketplace).
# network.target means interfaces configured; network-online.target means
# actually reachable. Pair with Wants= to pull it in.
Wants=network-online.target
After=network-online.target

# HACS API must be up before launching — chassis registers with the
# coordination server during bootstrap. Order matters at boot.
Wants=mcp-coordination.service
After=mcp-coordination.service

# Don't restart-loop on persistent errors: 3 attempts in 5 min, then give up.
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
User=root

# Working directory is the instance dir (where .mcp.json and .claude.json live)
WorkingDirectory=/mnt/coordinaton_mcp_data/instances/%i

# Launch reads --port from .hacs-identity automatically (no port arg needed)
ExecStart=/mnt/coordinaton_mcp_data/Human-Adjacent-Coordination/src/chassis/claude-code-channel/launch-claude-code-channel.sh --instance-id %i

# Land kills the tmux session as the instance user, frees the port
ExecStop=/mnt/coordinaton_mcp_data/Human-Adjacent-Coordination/src/chassis/claude-code-channel/land-claude-code-channel.sh --instance-id %i

# Launch script does its own health-check wait (up to 30s for /health to respond).
# Give it 60s margin total before systemd declares failure.
TimeoutStartSec=60
TimeoutStopSec=30

# Restart on failure (crash, OOM, etc.). Rate limits in [Unit] above.
Restart=on-failure
RestartSec=10

# Resource limits — claude-code can balloon during long sessions or large
# file operations. Cap memory so a runaway chassis doesn't OOM the host;
# raise file-descriptor limit because claude opens many fds.
MemoryMax=2G
LimitNOFILE=4096

# Hardening — instance dir is at /mnt/coordinaton_mcp_data/instances/%i/,
# not /home, so blocking /home costs nothing and tightens the surface.
# (ProtectSystem=strict + NoNewPrivileges would break sudo and npm; skip.)
ProtectHome=true

# Logging — journalctl -u claude-code-channel@<InstanceId>
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-channel-%i

[Install]
WantedBy=multi-user.target
