# src/chassis/openfang/openfang-approver@.service
# OpenFang Auto-Approver Sidecar Service Template
# Companion to openfang@.service — auto-approves tool executions
#
# Usage: Started automatically by openfang@.service (via Wants=)
#        Can also be managed independently:
#        systemctl start openfang-approver@Flair-2a84
#
# %i = instance ID (e.g., Flair-2a84)

[Unit]
Description=OpenFang Auto-Approver for HACS instance %i
BindsTo=openfang@%i.service
After=openfang@%i.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/mnt/coordinaton_mcp_data/instances/%i/openfang

# Read port from config.toml at startup
ExecStart=/bin/bash -c '\
  PORT=$(grep api_listen config.toml | grep -o "[0-9]\\+" | tail -1); \
  exec python3 auto-approve.py --port $PORT --interval 0.3'

Restart=on-failure
RestartSec=3

# Security — same as parent
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/mnt/coordinaton_mcp_data/instances/%i

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openfang-approver-%i

[Install]
WantedBy=openfang@%i.service
